# r3it/cert.target
[Unit]Description=r3it cert renewal timer
StopWhenUnneeded=yes