# r3it/cert.timer
[Unit]
Description=Cert renewal for r3it

[Timer]
OnBootSec=5min
OnCalendar=weekly

[Install]
WantedBy=timers.target